# patch_watcher/patch_watcher_daily.service
[Unit]
Description=Lustre Patch Watcher Daily Confidence Report
After=network-online.target

[Service]
Type=oneshot
User=patchwatcher
ExecStart=/shared/llm_code_and_review_tools/patch_watcher/daily_confidence.sh
TimeoutStartSec=120
Environment=HOME=/home/patchwatcher
Environment=PATH=/home/patchwatcher/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/shared/llm_code_and_review_tools/patch_watcher
